# packages/lakka/lakka_tools/cec_mini_kb/system.d/cec-mini-kb.service
[Unit]
Description=cec-mini-kb
After=retroarch.service
Requires=retroarch.service

[Service]
ExecStart=/usr/bin/cec-mini-kb --poweroff "shutdown -P now"
Restart=no

[Install]
WantedBy=retroarch.service
